# ddterm/dbus-com.github.amezin.ddterm.service.in
[Service]
Type=forking
BusName=com.github.amezin.ddterm
ExecStart="@LAUNCHER@" --launch-through-extension --attach-unit "%n"
Delegate=
