# test/etc/systemd/system/gnome-session-wayland@.service.d/mutter-highdpi.conf
[Service]
Environment="MUTTER_DEBUG_DUMMY_MODE_SPECS=1920x1080@30.0"
Environment="MUTTER_DEBUG_DUMMY_MONITOR_SCALES=2"
