# test/etc/systemd/system/gnome-session-wayland@.service.d/mutter-fractional.conf
[Service]
Environment="MUTTER_DEBUG_DUMMY_MODE_SPECS=2880x1620@30.0"
Environment="MUTTER_DEBUG_DUMMY_MONITOR_SCALES=1.5"
